# 3rdparty/libqhy-old/85-qhyccdtemplete.rules
# Udev rules for QHY CCD cameras
# QHY2,  QHY5, QHY6, QHY8, QHY10
#
# ---- Initial Release ----
# Clive Rogers.  2 May 2009
#
# Add rules for QHY2, QHY5, QHY6, QHY8
#
# ---- Update 1 ----
# YANG Xudong 21 Oct 2010
#
# Add rules for QHY10 and modified for libqhyccd
#
# ---- Update 2 ----
# Fixed rules for the already initialized devices
# and make them user accessibe,
#

# Camera inserted.
# Vendorid & Prodid: QHY VOTI

ACTION!="add", GOTO="qhy_vidpid_end"
SUBSYSTEM=="usb", ENV{DEVTYPE}=="usb_device", GOTO="qhy_vidpid_test"
SUBSYSTEM!="usb_device", GOTO="qhy_vidpid_end"

LABEL="qhy_vidpid_test"

ATTRS{idVendor}=="1618", ATTRS{idProduct}=="b618", RUN+="/sbin/fxload -t fx2 -I @FIRMWARE_INSTALL_DIR@/IMG0H.HEX -D $env{DEVNAME}"
ATTRS{idVendor}=="1618", ATTRS{idProduct}=="0412", RUN+="/sbin/fxload -t fx2 -I @FIRMWARE_INSTALL_DIR@/QHY2.HEX -D $env{DEVNAME}"
ATTRS{idVendor}=="1618", ATTRS{idProduct}=="0901", RUN+="/sbin/fxload -t fx2 -I @FIRMWARE_INSTALL_DIR@/QHY5LOADER.HEX -D $env{DEVNAME}"
SYSFS{idVendor}=="1618", SYSFS{idProduct}=="0901", RUN+="/bin/sleep 5"
ATTRS{idVendor}=="1618", ATTRS{idProduct}=="0901", RUN+="/sbin/fxload -t fx2 -I @FIRMWARE_INSTALL_DIR@/QHY5.HEX -D $env{DEVNAME}"
ATTRS{idVendor}=="1618", ATTRS{idProduct}=="1002", RUN+="/sbin/fxload -t fx2 -I @FIRMWARE_INSTALL_DIR@/QHY5LOADER.HEX -D $env{DEVNAME}"
SYSFS{idVendor}=="1618", SYSFS{idProduct}=="1002", RUN+="/bin/sleep 5"
ATTRS{idVendor}=="1618", ATTRS{idProduct}=="1002", RUN+="/sbin/fxload -t fx2 -I @FIRMWARE_INSTALL_DIR@/QHY5.HEX -D $env{DEVNAME}"
ATTRS{idVendor}=="0547", ATTRS{idProduct}=="1002", RUN+="/sbin/fxload -t fx2 -I @FIRMWARE_INSTALL_DIR@/QHY5LOADER.HEX -D $env{DEVNAME}"
SYSFS{idVendor}=="0547", SYSFS{idProduct}=="1002", RUN+="/bin/sleep 5"
ATTRS{idVendor}=="0547", ATTRS{idProduct}=="1002", RUN+="/sbin/fxload -t fx2 -I @FIRMWARE_INSTALL_DIR@/QHY5.HEX -D $env{DEVNAME}"
ATTRS{idVendor}=="1618", ATTRS{idProduct}=="0259", RUN+="/sbin/fxload -t fx2 -I @FIRMWARE_INSTALL_DIR@/QHY6.HEX -D $env{DEVNAME}"
ATTRS{idVendor}=="1618", ATTRS{idProduct}=="6000", RUN+="/sbin/fxload -t fx2 -I @FIRMWARE_INSTALL_DIR@/QHY8.HEX -D $env{DEVNAME}"
ATTRS{idVendor}=="1618", ATTRS{idProduct}=="6004", RUN+="/sbin/fxload -t fx2 -I @FIRMWARE_INSTALL_DIR@/QHY8L.HEX -D $env{DEVNAME}"
ATTRS{idVendor}=="1618", ATTRS{idProduct}=="1000", RUN+="/sbin/fxload -t fx2 -I @FIRMWARE_INSTALL_DIR@/QHY10.HEX -D $env{DEVNAME}"
ATTRS{idVendor}=="1618", ATTRS{idProduct}=="1200", RUN+="/sbin/fxload -t fx2 -I @FIRMWARE_INSTALL_DIR@/QHY12.HEX -D $env{DEVNAME}"
ATTRS{idVendor}=="1618", ATTRS{idProduct}=="2850", RUN+="/sbin/fxload -t fx2 -I @FIRMWARE_INSTALL_DIR@/IMG2P.HEX -D $env{DEVNAME}"
ATTRS{idVendor}=="1618", ATTRS{idProduct}=="b285", RUN+="/sbin/fxload -t fx2lp -I @FIRMWARE_INSTALL_DIR@/IMG2S.HEX -D $env{DEVNAME}"
ATTRS{idVendor}=="1618", ATTRS{idProduct}=="0920", RUN+="/sbin/fxload -t fx2 -I @FIRMWARE_INSTALL_DIR@/QHY5II.HEX -D $env{DEVNAME}"
ATTRS{idVendor}=="1618", ATTRS{idProduct}=="8300", RUN+="/sbin/fxload -t fx2 -I @FIRMWARE_INSTALL_DIR@/QHY9S.HEX -D $env{DEVNAME}"
ATTRS{idVendor}=="1618", ATTRS{idProduct}=="8310", RUN+="/sbin/fxload -t fx2 -I @FIRMWARE_INSTALL_DIR@/IC8300.HEX -D $env{DEVNAME}"
ATTRS{idVendor}=="1618", ATTRS{idProduct}=="6940", RUN+="/sbin/fxload -t fx2 -I @FIRMWARE_INSTALL_DIR@/QHY22.HEX -D $env{DEVNAME}"
ATTRS{idVendor}=="1618", ATTRS{idProduct}=="6740", RUN+="/sbin/fxload -t fx2 -I @FIRMWARE_INSTALL_DIR@/QHY21.HEX -D $env{DEVNAME}"
ATTRS{idVendor}=="1618", ATTRS{idProduct}=="8140", RUN+="/sbin/fxload -t fx2 -I @FIRMWARE_INSTALL_DIR@/QHY23.HEX -D $env{DEVNAME}"
ATTRS{idVendor}=="1618", ATTRS{idProduct}=="1110", RUN+="/sbin/fxload -t fx2 -I @FIRMWARE_INSTALL_DIR@/QHY11.HEX -D $env{DEVNAME}"
ATTRS{idVendor}=="04b4", ATTRS{idProduct}=="8613", RUN+="/sbin/fxload -t fx2 -I @FIRMWARE_INSTALL_DIR@/QHY16000.HEX -D $env{DEVNAME}"
ATTRS{idVendor}=="1618", ATTRS{idProduct}=="6002", RUN+="/sbin/fxload -t fx2 -I @FIRMWARE_INSTALL_DIR@/QHY8PRO.HEX -D $env{DEVNAME}"
LABEL="qhy_vidpid_end"

# This is the ID presented once the firmware has been loaded.
# idVendor and idProduct match Cypress FX2 USB chips.

ACTION!="add", GOTO="qhy-rules-end"
SUBSYSTEMS=="usb", GOTO="qhy-vipi"
GOTO="qhy-rules-end"

LABEL="qhy-vipi"

SUBSYSTEMS=="usb", ATTRS{idVendor}=="16c0", ATTRS{idProduct}=="a618", MODE="0666"
SUBSYSTEMS=="usb", ATTRS{idVendor}=="16c0", ATTRS{idProduct}=="081e", MODE="0666"
SUBSYSTEMS=="usb", ATTRS{idVendor}=="16c0", ATTRS{idProduct}=="296d", MODE="0666"
SUBSYSTEMS=="usb", ATTRS{idVendor}=="16c0", ATTRS{idProduct}=="081d", MODE="0666"
SUBSYSTEMS=="usb", ATTRS{idVendor}=="1618", ATTRS{idProduct}=="2972", MODE="0666"
SUBSYSTEMS=="usb", ATTRS{idVendor}=="1618", ATTRS{idProduct}=="1001", MODE="0666"
SUBSYSTEMS=="usb", ATTRS{idVendor}=="1618", ATTRS{idProduct}=="2851", MODE="0666"
SUBSYSTEMS=="usb", ATTRS{idVendor}=="1618", ATTRS{idProduct}=="a285", MODE="0666"
SUBSYSTEMS=="usb", ATTRS{idVendor}=="1618", ATTRS{idProduct}=="0921", MODE="0666"
SUBSYSTEMS=="usb", ATTRS{idVendor}=="1618", ATTRS{idProduct}=="8301", MODE="0666"
SUBSYSTEMS=="usb", ATTRS{idVendor}=="1618", ATTRS{idProduct}=="8311", MODE="0666"
SUBSYSTEMS=="usb", ATTRS{idVendor}=="1618", ATTRS{idProduct}=="6941", MODE="0666"
SUBSYSTEMS=="usb", ATTRS{idVendor}=="1618", ATTRS{idProduct}=="6741", MODE="0666"
SUBSYSTEMS=="usb", ATTRS{idVendor}=="1618", ATTRS{idProduct}=="8141", MODE="0666"
SUBSYSTEMS=="usb", ATTRS{idVendor}=="1618", ATTRS{idProduct}=="1111", MODE="0666"
SUBSYSTEMS=="usb", ATTRS{idVendor}=="1618", ATTRS{idProduct}=="1611", MODE="0666"
SUBSYSTEMS=="usb", ATTRS{idVendor}=="1618", ATTRS{idProduct}=="6003", MODE="0666"

LABEL="qhy-rules-end"
