# 3rdparty/libaltaircam/99-altaircam.rules
ACTION=="add", ATTR{idVendor}=="0547", RUN+="/bin/sh -c '/bin/echo 256 >/sys/module/usbcore/parameters/usbfs_memory_mb'"
ACTION=="add", ATTR{idVendor}=="16d0", RUN+="/bin/sh -c '/bin/echo 256 >/sys/module/usbcore/parameters/usbfs_memory_mb'"
# All Altair devices
SUBSYSTEMS=="usb", ATTR{idVendor}=="0547", MODE="0666"
SUBSYSTEMS=="usb", ATTR{idVendor}=="16d0", MODE="0666"
