# 3rdparty/libatik/99-atik.rules
ACTION=="add", ATTR{idVendor}=="03c3", RUN+="/bin/sh -c '/bin/echo 256 >/sys/module/usbcore/parameters/usbfs_memory_mb'"
ACTION=="add", ATTR{idVendor}=="20e7", RUN+="/bin/sh -c '/bin/echo 256 >/sys/module/usbcore/parameters/usbfs_memory_mb'"
ACTION=="add", ATTR{idVendor}=="04d8", RUN+="/bin/sh -c '/bin/echo 256 >/sys/module/usbcore/parameters/usbfs_memory_mb'"
ACTION=="add", ATTR{idVendor}=="1e10", RUN+="/bin/sh -c '/bin/echo 256 >/sys/module/usbcore/parameters/usbfs_memory_mb'"
# All Atik Cameras and filter wheels
SUBSYSTEMS=="usb", ATTR{idVendor}=="04b4", MODE="0666"
SUBSYSTEMS=="usb", ATTR{idVendor}=="20e7", MODE="0666"
SUBSYSTEMS=="usb", ATTR{idVendor}=="04d8", MODE="0666"
SUBSYSTEMS=="usb", ATTR{idVendor}=="1e10", MODE="0666"
