# 3rdparty/libsbig/51-sbig-debian.rules
# Increase USBFS memory
ACTION=="add", ATTR{idVendor}=="0d97", RUN+="/bin/sh -c '/bin/echo 256 >/sys/module/usbcore/parameters/usbfs_memory_mb'"

# SBIG ST-7/8/9/10/2K cameras
ACTION=="add", SUBSYSTEM=="usb", ATTR{idVendor}=="0d97", ATTR{idProduct}=="0001", \
RUN+="/sbin/fxload -I /lib/firmware/sbigucam.hex -D $env{DEVNAME}"

# SBIG 4K camera
ACTION=="add", SUBSYSTEM=="usb", ATTR{idVendor}=="0d97", ATTR{idProduct}=="0001", ATTRS{serial}=="3f980000.usb" \
RUN+="/sbin/fxload -I /lib/firmware/sbigpcam.hex -D $env{DEVNAME}"

# SBIG ST-L cameras
ACTION=="add", SUBSYSTEM=="usb", ATTR{idVendor}=="0d97", ATTR{idProduct}=="0002", \
RUN+="/sbin/fxload -I /lib/firmware/sbiglcam.hex -D $env{DEVNAME}"

# SBIG ST-402/1603/3200/8300 cameras
ACTION=="add", SUBSYSTEM=="usb", ATTR{idVendor}=="0d97", ATTR{idProduct}=="0003", \
RUN+="/sbin/fxload -I /lib/firmware/sbigfcam.hex -D $env{DEVNAME} -t fx2"

# Make all SBIG Cameras accessible post enumeration
SUBSYSTEM=="usb", ATTR{idVendor}=="0d97", MODE="0666"
